# test3/deploy/systemd/test3-frontend.service
[Unit]
Description=test3 frontend
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/srv/test3/frontend
Environment=NODE_ENV=production
EnvironmentFile=/srv/test3/frontend/.env.production
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
